# crontab
* * * * */5 ubuntu /bin/bash /home/ubuntu/apps/commencement/repository/run_on_server.sh fab $DEPLOYMENT_TARGET get_most_viewed >/var/log/commencement-crontab.log 2>&1
